# container/php/supervisor.ini
[supervisord]
nodaemon=true

[program:php-fpm]
command = php-fpm -F
user = root
autostart = true
